# systemd/gost-switchd.service
[Unit]
Description=Gost Auto Switch Daemon
After=network.target

[Service]
ExecStart=/usr/bin/python3 /opt/gost-switchd/gost-switchd.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
